# systemd/helix-frontend.service
[Unit]
Description=Helix Bot Builder Frontend (Production Server)
After=network.target

[Service]
Type=simple
User=talinoserver2
WorkingDirectory=/home/talinoserver2/Documents/bot-builder-suite
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=3000
ExecStartPre=/usr/bin/npm run build
ExecStart=/usr/bin/npm run preview -- --host 127.0.0.1 --port 3000
Restart=always
RestartSec=10
StandardOutput=append:/home/talinoserver2/Documents/bot-builder-suite/logs/frontend.log
StandardError=append:/home/talinoserver2/Documents/bot-builder-suite/logs/frontend.error.log

[Install]
WantedBy=multi-user.target
